# lab-rip-smpl/net/r5/etc/rc.local
#!/bin/sh -e

iptables -t nat -A POSTROUTING -o eth2 -j MASQUERADE

exit 0
